# etc/crontab
00 00 * * * python /app/manage.py get_latest_tracks